# redis/redisstart.sh
#!/bin/sh
#chkconfig: 2345 80 90
#description:redis_开机自启脚本

#开启redis服务
